# devops/app/supervisor/conf.d/aiohttp.conf
[program:aiohttp]
directory=/app
command=python3 main.py
user=app
autostart=true
autorestart=true
redirect_stderr=true
priority=800
